# home/.config/autostart/ssh-add.sh.desktop
[Desktop Entry]  
Exec=ssh-add  
Icon=dialog-scripts  
Name=ssh-add.sh  
Type=Application  
X-KDE-AutostartScript=true